# release/debian/bamboohr-slack-bot/DEBIAN/postrm
#!/bin/bash

set -e

systemctl disable bamboohr-slack-bot
systemctl daemon-reload
systemctl restart rsyslog
